# deploy/crontab
*/5 * * * * /var/www/flaviodesousa/deploy/flaviodesousa-crontab-twitter-poll.sh
